# data/etc/xdg/autostart/ayatana-webmail-autostart.desktop
[Desktop Entry]
Version=1.0
Name=Ayatana Webmail (service)
Comment=Display notifications about new mail
Icon=ayatanawebmail
Exec=ayatana-webmail
Type=Application
Categories=Email;Network;
